# src/systemd/certmonger.path.in
[Unit]
Description=Certificate monitoring and PKI enrollment
After=syslog.target network.target dbus.service

[Path]
DirectoryNotEmpty=@CM_STORE_REQUESTS_DIRECTORY@

[Install]
WantedBy=multi-user.target
